# hw/05/00-systemd-user.d/ns1-vde.service
[Unit]
Description=ns1 VDE
After=network.target

[Service]
Type=forking
ExecStart=vde_switch -sock /home/volpind/vde/sw1.sock -daemon

[Install]
WantedBy=default.target
